# fastapi/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn3 app:app -w 4 -k uvicorn.workers.UvicornWorker --bind unix:gunicorn.sock -m 777 --log-file=gunicornErrors.log
directory=/app/app
